# deploy/lelianahoy.service
[Unit]
Description=sabhoy.es (Next.js)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sabhoy.es
Environment=NODE_ENV=production
EnvironmentFile=-/opt/sabhoy.es/.env
# next start (puerto 3000; Nginx u otro proxy delante para 80/443)
ExecStart=/usr/bin/npm run start -- -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=10
# Mismo usuario que posee /opt/sabhoy.es (recomendado; si no, el servicio corre como root)
# User=tu-usuario-deploy
# Group=tu-usuario-deploy

[Install]
WantedBy=multi-user.target
