# deploy/sabhoy.service
[Unit]
Description=sabhoy.es (Next.js)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sabhoy.es
Environment=NODE_ENV=production
EnvironmentFile=-/opt/sabhoy.es/.env
# next start (puerto 3003; 3001 = sermestre, 3000 = lelianahoy, 3002 = beterahoy)
ExecStart=/usr/bin/npm run start -- -H 127.0.0.1 -p 3003
Restart=on-failure
RestartSec=10
# User=tu-usuario-deploy
# Group=tu-usuario-deploy

[Install]
WantedBy=multi-user.target
